# cron_file/127.0.0.1/localhost.cron
# simple example of crontab
* * * * * date >> /tmp/crontab_controller.test
