# cron_file/127.0.0.1/localhost.2.cron
# simple example of crontab
* * * * * uptime >> /tmp/crontab_controller.test
